# backup-image/crontab
# START CRON JOB
* * * * * python /app/backup.py
# END CRON JOB